# systemd/assistant-deploy.service
[Unit]
Description=Auto-deploy telegram-assistant from GitHub
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/home/assistant/telegram-assistant/systemd/deploy.sh
